# deploy/docker/supervisord.conf
# OpsKG 单镜像进程管理（S13-4 多阶段构建）
#
# supervisord 同时管理 nginx + uvicorn 两个进程：
#   - nginx：前端静态资源 + 反向代理
#   - uvicorn：后端 FastAPI
#
# 任一进程异常退出 → supervisord 自动重启该进程
# 容器仅在 supervisor 自身被 SIGKILL 时退出（K8s/Docker 会重启 Pod）
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid
# P1-2: 以 opskg 非 root 用户运行 supervisord
user=opskg

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# ── nginx：前端 + 反向代理 ──
# P1-2: nginx 监听 8080（非特权），以 opskg 用户运行
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=opskg
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=10
priority=10
stdout_logfile=/var/log/nginx/supervisor-stdout.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/nginx/supervisor-stderr.log
stderr_logfile_maxbytes=5MB

# ── uvicorn：后端 FastAPI ──
# --workers 2 默认值；可通过 OPSKG_UVICORN_WORKERS 环境变量调整
# 注：supervisord 不直接支持环境变量插值，需通过 entrypoint 包装
# 以 opskg 用户运行（降权，安全最佳实践）
[program:uvicorn]
command=/usr/local/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/app
user=opskg
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/app",PYTHONUNBUFFERED="1"
